# realestatefetcher/crontab/root
*/15 * * * * node /home/node/app/index.js
0 3 * * * node /home/node/app/index.js --full-sync